# supervisord.conf
[supervisord]
nodaemon=true

[program:kumquat-frontend]
directory=/kumquat/admin
command=nodemon server.js
